# scripts/com.jarvis.daily_restart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jarvis.daily_restart</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>~/jarvis/scripts/daily_restart.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>0</integer>
    <key>Minute</key>
    <integer>5</integer>
  </dict>
  <key>KeepAlive</key>
  <false/>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>~/jarvis/logs/daily_restart.log</string>
  <key>StandardErrorPath</key>
  <string>~/jarvis/logs/daily_restart.log</string>
</dict>
</plist>
